# core/config/supservisor.conf
[program:postfix]
command=/postfix.sh
autostart=true
autorestart=true
user=root

[program:dovecot]
command=/usr/sbin/dovecot -c /etc/dovecot/dovecot.conf -F
autostart=true
autorestart=true
user=root
